# pigen-splitstreamer/01-copy-server-binaries/splitstreamer.service
[Unit]
Description=SplitStreamer
After=signaling.target
Wants=signaling.target

[Service]
WorkingDirectory=/boot
ExecStart=/usr/bin/server
Restart=true

[Install]
WantedBy=multi-user.target
